# com.soca.weather-bot.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Example launchd plist for soca-weather-bot.

  Setup:
    1. 複製這份檔案到 ~/Library/LaunchAgents/com.soca.weather-bot.plist
    2. 把所有 /Users/YOUR_USERNAME/ 改成你的實際路徑
    3. launchctl load ~/Library/LaunchAgents/com.soca.weather-bot.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.soca.weather-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/YOUR_USERNAME/projects/soca-weather-bot/weather_bot.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/projects/soca-weather-bot/logs/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/projects/soca-weather-bot/logs/launchd.err.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>
</dict>
</plist>
